# uesp-update-lsearch.cron
#!/bin/sh

cd /home/uesp/lucenesearch
./update
